# usbhub-wakeup/etc/udev/rules.d/10-wakeup.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{bDeviceClass}=="09", RUN+="/usr/local/bin/usbhub-wakeup"
